# files/qemu-aarch64/balance-irq/etc/init.d/balance-irq
#!/bin/sh /etc/rc.common

START=42
 
start() {        
	/usr/sbin/gen-balance-irq-conf.sh
	/usr/sbin/balethirq.pl
}                 
